# files/system/etc/skel/.config/autostart/desktop-setup.desktop
[Desktop Entry]
Type=Application
Name=Setup
Exec=$HOME/.local/bin/desktop-setup.sh
Icon=utilities-terminal
Terminal=true
Categories=Utility;