# files/systemd/system/skel-init.service
[Unit]
Description=Skel initialization for existing users (bootc-managed)
After=systemd-user-sessions.service
Before=greetd.service
ConditionPathExists=!/var/lib/dank-blue-build/skel-init-done

[Service]
Type=oneshot
ExecStart=/usr/libexec/dank-blue-build/skel-init

[Install]
WantedBy=multi-user.target